# src/octane-nginx/etc/s6-overlay/s6-rc.d/octane/run
#!/command/execlineb -P
with-contenv
s6-notifyoncheck -d
/usr/bin/php /var/www/html/artisan octane:start
